# etc/service/targets/run
#!/bin/sh

exec 2>&1
exec python /usr/bin/prom_targets.py