# etc/service/admin/run
#!/bin/sh

exec 2>&1
exec python /usr/bin/admin.py